# systemd/clean-motion-videos.service
[Unit]
Description=Remove specific videos used by Motion

[Service]
Type=simple
ExecStart=/usr/local/bin/clean-motion-videos -n
