# systemd/clean-motion-videos.timer
[Unit]
Description=Run clean-motion-videos daily

[Timer]
OnCalendar=daily 

[Install]
WantedBy=timers.target
